# core/imageroot/etc/systemd/system/user@.service.d/20-after-redis.conf
[Unit]
#
# Start the session of rootless modules after Redis
#
After=redis.service
