# core/imageroot/etc/systemd/system/user@.service.d/30-fix-xdg-state.conf
[Unit]
After=fix-xdg-state@%i.service
Wants=fix-xdg-state@%i.service
