# mail/imageroot/systemd/user/dovecot.service
[Unit]
Description=Podman dovecot.service
Documentation=man:podman-generate-systemd(1)
BindsTo=mail.service
After=mail.service

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
EnvironmentFile=%S/state/environment
Restart=always
TimeoutStopSec=70
ExecStartPre=/bin/rm -f %t/dovecot.pid %t/dovecot.ctr-id
ExecStart=/usr/bin/podman run \
    --detach \
    --conmon-pidfile=%t/dovecot.pid \
    --cidfile=%t/dovecot.ctr-id \
    --cgroups=no-conmon \
    --pod-id-file=%t/mail.pod-id \
    --log-opt=tag=%u \
    --replace --name=%N \
    --volume=vmail:/var/lib/vmail:Z \
    --env-file=./environment \
    ${DOVECOT_IMAGE}

ExecStop=/usr/bin/podman stop --ignore --cidfile %t/dovecot.ctr-id -t 10
ExecStopPost=/usr/bin/podman rm --ignore -f --cidfile %t/dovecot.ctr-id
PIDFile=%t/dovecot.pid
Type=forking
WorkingDirectory=%S/state
