# core/imageroot/etc/systemd/system/dnf-automatic-notifyonly.service.d/10-ns8.conf
[Service]
# An active NS8 subscription is not compatible with dnf-automatic
ExecCondition=runagent -m node check-unmanaged-updates
