# core/imageroot/etc/systemd/system/wg-quick@wg0.service.d/10-ns8.conf
[Unit]
ConditionPathExists=/etc/wireguard/wg0.conf
