# core/imageroot/etc/systemd/system/apply-updates.service
[Unit]
Description=Download and install core, modules and OS updates

[Service]
Type=oneshot
ExecStart=runagent apply-updates
SyslogIdentifier=%N
