# core/imageroot/etc/systemd/system/overlay-cleanup.service
[Unit]
Description=Cleanup the overlay-merged volume

[Service]
RemainAfterExit=true
ExecStart=/usr/local/sbin/overlay-cleanup
Type=oneshot
