# core/imageroot/etc/systemd/system/password-warning.timer
[Unit]
Description=Domain user password notification trigger

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
