# core/imageroot/etc/systemd/system/password-warning.service
[Unit]
Description=Domain user password notification

[Service]
Type=oneshot
ExecStart=runagent -m cluster notify-password-warning
