# core/imageroot/etc/systemd/system/phonehome.timer
[Unit]
Description=PhoneHome Trigger

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=21600

[Install]
WantedBy=timers.target
