# core/imageroot/etc/systemd/system/phonehome.service
[Unit]
Description=PhoneHome

[Service]
Type=oneshot
ExecStart=/usr/local/bin/phonehome
