# core/imageroot/etc/systemd/system/phonehome.service
[Unit]
Description=PhoneHome

[Service]
Type=oneshot
ExecStart=runagent send-phonehome
