# core/imageroot/etc/systemd/system/rclone-webdav.service
[Unit]
Description=Rclone WebDAV server

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Environment=OPTIONS=
Environment=BACKUP_VOLUME=rclone-webdav
EnvironmentFile=/etc/nethserver/core.env
EnvironmentFile=-/var/lib/nethserver/node/state/rclone-webdav.env
Restart=always
TimeoutStopSec=120
ExecStartPre=/bin/rm -f %t/%N.pid %t/%N.cid
ExecStart=/usr/bin/podman run \
    --conmon-pidfile=%t/%N.pid \
    --cidfile=%t/%N.cid \
    --cgroups=no-conmon \
    --detach \
    --log-opt=tag=%N \
    --replace --name=%N \
    --network=host \
    --user=restic:restic \
    --volume=${BACKUP_VOLUME}:/srv/repo:z \
    --entrypoint=[] \
    ${RESTIC_IMAGE} rclone serve webdav /srv/repo --addr :4694 $OPTIONS

ExecStop=/usr/bin/podman stop --ignore --cidfile %t/%N.cid -t 90
ExecStopPost=/usr/bin/podman rm --ignore -f --cidfile %t/%N.cid
PIDFile=%t/%N.pid
Type=forking
SuccessExitStatus=143

[Install]
WantedBy=default.target
