# core/imageroot/etc/systemd/system/send-inventory.service
[Unit]
Description=Subscription inventory

[Service]
Type=oneshot
ExecStart=runagent send-inventory
SyslogIdentifier=%N
