# core/imageroot/etc/systemd/system/send-inventory.service
[Unit]
Description=Subscription inventory

[Service]
Type=simple
ExecStart=runagent send-inventory
SyslogIdentifier=send-inventory
