# core/imageroot/etc/systemd/system/send-backup.service
[Unit]
Description=Subscription cluster cloud backup
ConditionPathExists=/var/lib/nethserver/cluster/state/backup/passphrase

[Service]
Type=oneshot
ExecStart=runagent cluster-backup
ExecStart=runagent send-cluster-backup
SyslogIdentifier=%N
