# traefik/imageroot/systemd/user/certificate-exporter.path
[Unit]
Description=Monitor acme.json file for changes

[Path]
PathChanged=%h/.local/share/containers/storage/volumes/traefik-acme/_data/acme.json
Unit=certificate-exporter.service
