# core/imageroot/etc/profile.d/nethserver.sh
#
# This is useful during development. We should remove it later.
#

alias logcli="/var/lib/nethserver/node/logcli.sh"
