# core/imageroot/etc/profile.d/nethserver.sh

#
# Copyright (C) 2025 Nethesis S.r.l.
# SPDX-License-Identifier: GPL-3.0-or-later
#

if [ -z "$BASH" ]; then
    return
fi

export HISTSIZE=2000
export HISTFILESIZE=2000
export HISTTIMEFORMAT="%F %T "
